# scripts/smartserow-ui sample.service
[Unit]
Description=Smart Serow UI
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/smartserow/bundle

# GBM backend is compiled in. Just specify bundle path.
ExecStart=/opt/smartserow/bundle/smartserow_ui --bundle=/opt/smartserow/bundle

# Restart on crash
Restart=always
RestartSec=3

# Environment for DRM/KMS access
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=HOME=/home/pi
Environment=LD_LIBRARY_PATH=/opt/smartserow/bundle/lib

# Give time for GPU to initialize on boot
ExecStartPre=/bin/sleep 2

# Ensure clean shutdown
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
